# elk/etc/supervisor/conf.d/kibana.conf
[program:kibana]
command=/opt/kibana/bin/kibana -c /etc/kibana/kibana.yml -p 80 -l /var/log/kibana.log
